# rpmbuild/SOURCES/snap-firefox-5187.mount
[Unit]
Description=Mount unit for firefox, revision 5187
After=snapd.mounts-pre.target
Before=snapd.mounts.target

[Mount]
What=/var/lib/snapd/snaps/firefox_5187.snap
Where=/snap/firefox/5187
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=snapd.mounts.target
WantedBy=multi-user.target
